# rpi/tfluna/tfluna.service.default
[Unit]
Description=TF Luna Systemd Service Unit
After=network.target

[Service]
Type=simple
User=USER
ExecStart=/usr/bin/python3 PATH
Restart=on-failure

[Install]
WantedBy=multi-user.target

